# pkg/systemd/baseboxd.service.in
[Unit]
Description=baseboxd
After=network.target

[Service]
Type=simple
ExecStart=@bindir@/baseboxd --tryfromenv=port
EnvironmentFile=-@sysconfdir@/sysconfig/baseboxd
Restart=always

[Install]
WantedBy=multi-user.target
