# pkg/systemd/baseboxd.service.in
[Unit]
Description=baseboxd
After=network.target

[Service]
Type=notify
ExecStart=@bindir@/baseboxd
ExecStopPost=@bindir@/baseboxd-knet-reset.py
EnvironmentFile=-@sysconfdir@/sysconfig/baseboxd
Restart=always
NotifyAccess=main

[Install]
WantedBy=multi-user.target
