# pkg/systemd/baseboxd.service.in
[Unit]
Description=baseboxd
After=network.target

[Service]
Type=notify
ExecStart=@bindir@/baseboxd
EnvironmentFile=-@sysconfdir@/sysconfig/baseboxd
Restart=always
NotifyAcess=main

[Install]
WantedBy=multi-user.target
